# deploy/comment-tracker.service
[Unit]
Description=Comment Tracker Service
After=network.target

[Service]
Type=simple
User=callfans
WorkingDirectory=/home/callfans/data/comment-tracker
Environment=TZ=Asia/Shanghai
ExecStart=/home/callfans/miniconda3/bin/conda run -n py312 python /home/callfans/data/comment-tracker/main.py
Restart=always
RestartSec=10
StandardOutput=append:/home/callfans/data/comment-tracker/logs/service.log
StandardError=append:/home/callfans/data/comment-tracker/logs/error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
